# tom-cat/tomcat-7.0.61/bin/service
#!/bin/sh

### BEGIN INIT INFO
# Provides:          openkm
# Required-Start:    $remote_fs $syslog
# Required-Stop:     $remote_fs $syslog
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Short-Description: Start and stop OpenKM
# Description:       Enable OpenKM service provided by daemon.
### END INIT INFO

ECHO=/bin/echo
TEST=/usr/bin/test
TOMCAT_USER=openkm
TOMCAT_HOME=/home/openkm/tomcat
TOMCAT_START_SCRIPT=$TOMCAT_HOME/bin/startup.sh
TOMCAT_STOP_SCRIPT=$TOMCAT_HOME/bin/shutdown.sh

$TEST -x $TOMCAT_START_SCRIPT || exit 0
$TEST -x $TOMCAT_STOP_SCRIPT || exit 0

start() {
    # Wine stuff
    #$ECHO "Starting Xvfb"
    #Xvfb :1 &
    #export DISPLAY=:1

    $ECHO -n "Starting OpenKM"
    su - $TOMCAT_USER -c "$TOMCAT_START_SCRIPT &"
    $ECHO "."
}

stop() {
    # Wine stuff
    #$ECHO "Stopping Xvfb"
    #killall Xvfb

    $ECHO -n "Stopping OpenKM"
    su - $TOMCAT_USER -c "$TOMCAT_STOP_SCRIPT 60 -force &"
    while [ "$(ps -fu $TOMCAT_USER | grep java | grep tomcat | wc -l)" -gt "0" ]; do
        sleep 5; $ECHO -n "."
    done
    $ECHO "."
}

case "$1" in
    start)
        start
        ;;
    stop)
        stop
        ;;
    restart)
        stop
        sleep 30
        start
        ;;
    *)
        $ECHO "Usage: openkm {start|stop|restart}"
        exit 1
esac
exit 0

